# simple-mcp.service
[Unit]
Description=Simple MCP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/simple-mcp --config /etc/simple-mcp/simple-mcp.yaml --listen-addr :8080
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
